# bin/systemd/gunicorn-uat.socket
#!/bin/sh

[Unit]
Description=gunicorn uat socket

[Socket]
ListenStream=/run/gunicorn-uat.sock

[Install]
WantedBy=sockets.target
